# root/etc/s6-overlay/s6-rc.d/main/run
#!/usr/bin/with-contenv sh

UMASK_SET=${UMASK_SET:-022}
umask "$UMASK_SET"

exec s6-setuidgid abc python3 /abs_library_tools/src/alt/daemons.py
